# deploy/systemd/pradyos-admission.service
[Unit]
Description=PRADY OS — Admission Bridge (constitutional gate for ORACLE proposals, Plane 5)
Documentation=file:///opt/pradyos/docs/api/oracle.md
# The bridge subscribes to oracle.proposal events from the ORACLE daemon
# and must start after both ORACLE and IMPERIUM are live.
After=pradyos-oracle.service pradyos-imperium.service
Requires=pradyos-imperium.service
# Restart alongside ORACLE if it cycles.
BindsTo=pradyos-oracle.service

[Service]
Type=exec
ExecStart=/opt/pradyos/.venv/bin/pradyos-admission
Restart=on-failure
RestartSec=3
StartLimitBurst=5
StartLimitIntervalSec=30

User=pradyos
Group=pradyos

# ---------- Runtime configuration ----------
Environment=PRADYOS_AUDIT_PATH=/var/log/pradyos/audit.jsonl
Environment=PRADYOS_STATE_PATH=/var/lib/pradyos/state
Environment=PRADYOS_BUS_BACKEND=redis
Environment=PRADYOS_REDIS_URL=redis://127.0.0.1:6379/0
EnvironmentFile=-/etc/pradyos/admission.env

# ---------- Persistent directories ----------
LogsDirectory=pradyos
StateDirectory=pradyos

# ---------- Hardening (Python-ML-safe) ----------
# The bridge imports the ORACLE package (numpy / scikit-learn stack), so
# MemoryDenyWriteExecute= and SystemCallFilter= are deliberately absent —
# they kill the ML imports at load time.
ProtectSystem=strict
ReadWritePaths=/var/log/pradyos /var/lib/pradyos
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
NoNewPrivileges=true
ProtectKernelModules=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
LockPersonality=true
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
